# ui/telcoin-ui.service
[Unit]
Description=Telcoin Network Node Manager UI
After=network.target

[Service]
Type=simple
User=telcoin-ui
WorkingDirectory=/opt/telcoin-ui
ExecStart=/usr/bin/python3 /opt/telcoin-ui/server.py
Restart=on-failure
RestartSec=5

# Hardening. NoNewPrivileges is intentionally NOT set: the service shells out
# to `sudo ...` via the narrow /etc/sudoers.d/telcoin-ui whitelist, and
# NoNewPrivileges would block sudo's setuid transition. The attack surface is
# bounded instead by the localhost-only bind (server.py listens on 127.0.0.1)
# plus the whitelisted helper commands.
#
# ProtectSystem and ProtectHome apply to the WHOLE service mount namespace --
# including the root helpers we sudo to. Those helpers legitimately write /etc
# (groupadd/useradd during setup, the systemd unit, /etc/telcoin, ufw rules) and
# /home (the service user's reth cache).
#
# ProtectSystem=full was tried with ReadWritePaths=/etc /home to carve those trees
# back, but on this kernel/systemd version ProtectSystem=full still bind-mounts
# /etc read-only (ro,relatime) regardless of ReadWritePaths -- so groupadd failed
# with "cannot lock /etc/group" and node removal could not delete the unit. Set
# ProtectSystem=false so /etc is writable in the namespace. This does not reduce
# security in practice: the unprivileged Flask process (User=telcoin-ui) still
# cannot write /etc by ordinary Unix permissions; only the root helper invoked
# via sudo can, which is exactly the intent.
ProtectHome=read-only
ProtectSystem=false
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
